# setup/bridge-bot.service
[Unit]
Description=Bot for load testing InterBTC
After=network.target

[Service]
EnvironmentFile=/opt/interbtc/bridge-bot/.env.testnet
Type=simple
ExecStart=/usr/local/bin/node /opt/interbtc/bridge-bot/build/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target